# customization/rusty-fedora/profile.d/rusty-aliases.sh
alias ls='eza'
alias cat='bat --style="header,grid"'
alias top='btm --battery'
alias du='dust -p'
